# dockerfile-all-in-one/conf/supervisord.conf
[supervisord]
nodaemon=true
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


# ── Infrastructure (priority 1-5, start first)

[program:mariadb]
command=/usr/sbin/mariadbd --user=mysql --skip-name-resolve --character-set-server=utf8mb4 --collation-server=utf8mb4_unicode_ci --skip-innodb-read-only-compressed
user=root
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis-cache]
command=redis-server --port 6379 --bind 127.0.0.1
user=redis
autostart=true
autorestart=true
priority=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis-queue]
command=redis-server --port 6380 --bind 127.0.0.1
user=redis
autostart=true
autorestart=true
priority=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── nginx (priority 10)

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── frappe processes (priority 20, higher startretries for DB/Redis readiness) ─

[program:frappe-web]
command=/home/frappe/frappe-bench/env/bin/gunicorn --chdir=/home/frappe/frappe-bench/sites --bind=0.0.0.0:8000 --threads=4 --workers=2 --worker-class=gthread --worker-tmp-dir=/dev/shm --timeout=120 --preload frappe.app:application
directory=/home/frappe/frappe-bench
user=frappe
environment=HOME="/home/frappe"
autostart=true
autorestart=true
startretries=10
startsecs=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:frappe-socketio]
command=/home/frappe/.nvm/versions/node/v24.12.0/bin/node /home/frappe/frappe-bench/apps/frappe/socketio.js
directory=/home/frappe/frappe-bench
user=frappe
environment=HOME="/home/frappe"
autostart=true
autorestart=true
startretries=10
startsecs=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:frappe-worker-short]
command=/usr/local/bin/bench worker --queue short,default
directory=/home/frappe/frappe-bench
user=frappe
environment=HOME="/home/frappe"
autostart=true
autorestart=true
startretries=10
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:frappe-worker-long]
command=/usr/local/bin/bench worker --queue long,default,short
directory=/home/frappe/frappe-bench
user=frappe
environment=HOME="/home/frappe"
autostart=true
autorestart=true
startretries=10
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:frappe-scheduler]
command=/usr/local/bin/bench schedule
directory=/home/frappe/frappe-bench
user=frappe
environment=HOME="/home/frappe"
autostart=true
autorestart=true
startretries=10
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
